# crontab
15 12-13 * * * python3.8 main.py >> /logfile
