# crontab
53 12-15 * * * python3.8 main.py >> /logfile
